# systemd/ruroco.service
[Unit]
Description=ruroco service
Requires=network-online.target ruroco.socket ruroco-commander.service
After=network-online.target ruroco.socket ruroco-commander.service

[Service]
Type=simple
ExecStart=/usr/local/bin/ruroco-server --config /etc/ruroco/config.toml
Restart=always
User=ruroco
Group=ruroco

RemoveIPC=true
# No capabilities are needed: the UDP port is bound by systemd (ruroco.socket) and the fd is inherited
# via socket activation, so the service never calls bind() in this setup. (socket.rs has a bind()
# fallback, but it only runs without socket activation, which the shipped units never trigger.)
CapabilityBoundingSet=
AmbientCapabilities=
DeviceAllow=
KeyringMode=private
LockPersonality=true
MemoryDenyWriteExecute=true
NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
PrivateUsers=true
ProcSubset=pid
ProtectClock=true
ProtectControlGroups=true
ProtectHome=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectProc=noaccess
ProtectSystem=strict
# Mutable state (blocklist.msgpck) lives in StateDirectory (/var/lib/ruroco), pointed at by
# blocklist_dir in config.toml. With no ReadWritePaths, /etc/ruroco is fully read-only, so a
# compromised server cannot overwrite keys or config — only its own counter state.
StateDirectory=ruroco
# The server only inherits its UDP fd and connects (AF_UNIX) to the commander; it never bind()s
# under socket activation, so deny every bind() attempt.
SocketBindDeny=any
# AF_UNIX is correct ONLY under socket activation (inherited UDP fd + the AF_UNIX connect to the
# commander). If you stop using ruroco.socket and let socket.rs bind [::]:34020 itself, you must
# add AF_INET AF_INET6 here and relax SocketBindDeny, or the server cannot open its socket.
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
SystemCallErrorNumber=EPERM
SystemCallFilter=@system-service
SystemCallFilter=~@chown
SystemCallFilter=~@clock
SystemCallFilter=~@cpu-emulation
SystemCallFilter=~@debug
SystemCallFilter=~@keyring
SystemCallFilter=~@memlock
SystemCallFilter=~@module
SystemCallFilter=~@mount
SystemCallFilter=~@obsolete
SystemCallFilter=~@privileged
SystemCallFilter=~@raw-io
SystemCallFilter=~@reboot
SystemCallFilter=~@resources
SystemCallFilter=~@setuid
UMask=0077

[Install]
WantedBy=multi-user.target
